# tiles/experiment-launcher.service
[Unit]
Description=Config-driven script launcher
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/pi/tile-management/tiles/experiment-launcher.py /home/pi/tile-management/tiles/experiment-config.yaml
WorkingDirectory=/home/pi/tile-management
Restart=on-failure
KillMode=process
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
